# home-configuration/ubuntu-server/tor/etc/systemd/system/dnscrypt-proxy.socket.d/override.conf
[Socket]
ListenStream=127.0.0.1:9153
ListenDatagram=127.0.0.1:9153
